# backup.timer
[Unit]
Description=Backup /opt/data Timer

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true
Unit=backup.service

[Install]
WantedBy=timers.target
